# src/deb/control/postinst
#!/bin/sh

echo "Installing init script"
sudo update-rc.d game-of-three-server defaults