# src/deb/control/prerm
#!/bin/sh

echo "Uninstalling init script"
sudo update-rc.d -f game-of-three-server remove